# shell-roboshop-v2/catalogue.service
[Unit]
Description = $app_name Service

[Service]
User=roboshop
Environment=MONGO=true
// highlight-start
Environment=MONGO_URL="mongodb://mongodb.devopswithrani.online:27017/$app_name"
// highlight-end
ExecStart=/bin/node /app/server.js
SyslogIdentifier=$app_name

[Install]
WantedBy=multi-user.target